# config/systemd/user/activity-logger.service
[Unit]
Description=Hyprland activity logger — stream events to Telegram
# Must run inside an active Hyprland session (socket lives there)
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
# Give Hyprland 3s to write its socket before we try to connect
ExecStartPre=/usr/bin/sleep 3
ExecStart=/bin/bash %h/.config/hypr/scripts/activity-logger.sh
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
