# config/systemd/user/boot-report.service
[Unit]
Description=Send boot error/warning report to Telegram
# Run after network is up so curl can reach api.telegram.org
Wants=network-online.target
After=network-online.target
# Also wait for graphical session so journald has collected startup logs
After=graphical-session.target

[Service]
Type=oneshot
# Wait 30s for systemd to finish starting all units before sampling failed ones
ExecStartPre=/usr/bin/sleep 30
ExecStart=/bin/bash %h/.config/hypr/scripts/boot-report.sh --boot
StandardOutput=journal
StandardError=journal
# Non-fatal — never block login
RemainAfterExit=no

[Install]
WantedBy=graphical-session.target
